# src/fpm/rootfs/etc/services.d/syslog-stdout/run
#!/bin/sh
set -e

s6-setuidgid root syslog-stdout
