# src/fpm/rootfs/etc/services.d/syslogd/run
#!/bin/sh
set -e -u

exec s6-setuidgid root syslogd -O /dev/stderr -S -n
