# src/fpm/rootfs/etc/services.d/php-fpm/run
#!/bin/sh
set -e -u

s6-svwait /run/s6/services/syslogd

exec s6-setuidgid root php-fpm
